# deploy/systemd/link-codex-to-discord.service.example
[Unit]
Description=Link Codex To Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/link-codex-to-discord
Environment=CTI_HOME=/opt/link-codex-to-discord/.runtime-home
ExecStart=/usr/bin/env node /opt/link-codex-to-discord/dist/daemon.mjs
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
